# deploy/gunicorn.service
[Unit]
Description=Gunicorn instance to serve DocuExpress
After=network.target

[Service]
User=vladtrix
Group=www-data
WorkingDirectory=/home/vladtrix/DOCUEXPRESS PAGINA
# Use the project's virtualenv (adjust if your venv path is different)
Environment=PATH=/home/vladtrix/DOCUEXPRESS PAGINA/.venv/bin
Environment=RATELIMIT_ENABLED=False
# ExecStart: use the gunicorn from the venv and the WSGI application
ExecStart=/home/vladtrix/DOCUEXPRESS PAGINA/.venv/bin/gunicorn --workers 3 --bind unix:/run/gunicorn-docuexpress.sock wsgi:application

[Install]
WantedBy=multi-user.target
