# deploy/gunicorn.service
[Unit]
Description=Gunicorn instance to serve DocuExpress
After=network.target

[Service]
User=vladtrix
Group=www-data
# The project path contains a space; to avoid systemd parsing issues we use a shell wrapper
# ExecStart will cd into the project dir and exec the venv gunicorn.
Environment=RATELIMIT_ENABLED=False
ExecStart=/bin/bash -c 'cd "/home/vladtrix/DOCUEXPRESS PAGINA" && exec "/home/vladtrix/DOCUEXPRESS PAGINA/.venv/bin/gunicorn" --workers 3 --bind unix:/run/gunicorn-docuexpress.sock wsgi:application'

[Install]
WantedBy=multi-user.target
